# deploy/mynas-storage-recover.service
[Unit]
Description=Recover MyNAS NAS volume mount
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/mynas-storage-recover
